# services/linux/ai-monitor.service
[Unit]
Description=AI Monitor - Claude Code Session Monitor
After=network.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/run.sh
Restart=on-failure
RestartSec=5
Environment=AI_MONITOR_PORT=__PORT__
StandardOutput=append:__INSTALL_DIR__/logs/ai-monitor.log
StandardError=append:__INSTALL_DIR__/logs/ai-monitor.err

[Install]
WantedBy=default.target
